# examples/coreelec/filesystem/usr/lib/systemd/system/khadas-fan.service
[Unit]
Description=Khadas-fan

[Service]
Type=simple
ExecStart=/opt/fan monitor

[Install]
WantedBy=basic.target
